# tweet_weather.cron
* 6 * * * python tweet_weather.py
* 9 * * * python tweet_weather.py
* 12 * * * python tweet_weather.py
* 15 * * * python tweet_weather.py
* 16 * * * python tweet_weather.py
* 18 * * * python tweet_weather.py
* 21 * * * python tweet_weather.py
